# nginx_cgi_python_json_server/supervisord.conf
[program:myapp]
directory=/app
command=gunicorn --workers 3 --bind unix:/var/www/application.sock wsgi:app

[program:nginx]
command=/docker-entrypoint.sh nginx
